# ops/pi-synthetic-monitor/zadiag-synthetic-monitor.service
[Unit]
Description=Zadiag production synthetic Web Push monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/idriss/dev/zadiag
EnvironmentFile=/home/idriss/dev/zadiag/.pi-monitor/env
ExecStart=/usr/bin/xvfb-run -a -s "-screen 0 1280x800x24" /usr/bin/node /home/idriss/dev/zadiag/ops/pi-synthetic-monitor/supervisor.mjs
Restart=always
RestartSec=10
TimeoutStopSec=20

[Install]
WantedBy=default.target
